# config/spawn-fcgi.conf
[program:spawn-fcgi]
command = spawn-fcgi -p 9002 -n -- /usr/sbin/fcgiwrap
user = ubuntu
stdout_logfile = /var/log/supervisor/spawn-fcgi.log

